# root/etc/services.d/librespot/run
#!/usr/bin/execlineb -P
librespot --backend pipe --device /music/spotify -n forked-daapd --cache /tmp
