# root/etc/services.d/dbus/run
#!/usr/bin/with-contenv bash

exec \
	dbus-daemon --system --nofork


