# root/etc/s6-overlay/s6-rc.d/svc-dbus/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec \
    dbus-daemon --system --nofork
