# root/etc/services.d/forked/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/sbin/owntone -f \
	-P /daapd-pidfolder/owntone.pid
