# root/etc/services.d/forked/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc /app/sbin/forked-daapd -f -P /daapd-pidfolder/forked-daapd.pid

