# root/etc/services.d/avahi/run
#!/usr/bin/with-contenv bash

until [[ -e /var/run/dbus/system_bus_socket ]]; do
sleep 1s
done

exec \
	avahi-daemon --no-chroot
