# root/etc/s6-overlay/s6-rc.d/svc-librespot/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

librespot --backend pipe --device /music/spotify -n "${SPOTIFY_CONNECT_NAME:OwnTone}" --cache /tmp
