# root/etc/s6-overlay/s6-rc.d/svc-librespot/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

eval librespot --backend pipe --device /music/spotify --cache /tmp ${LIBRESPOT_ARGS:-"-n OwnTone"}
